# templates/init-httpd.conf.erb
# This file is managed by Puppet.
# Attemptes to modify it will be overwritten at the next puppet run.
# httpd-<%= @shortname -%> - Apache HTTP Server for <%= @instance %>
#

# description     "Apache HTTP Server for <%= @instance -%>"

start on virtual-filesystems
stop on runlevel [06]

respawn
respawn limit 5 30

<% unless @enabled -%>
manual
<% end -%>

setuid <%= @instance %>
setgid <%= @instance %>
console log

script
	exec  <%= @prefix -%>/bin/httpd -f <%= @vhostconf -%>/<%= @instance -%>/httpd.conf -DNO_DETACH -k start
end script
